# doc/dev-scripts/myLogger.service
[Unit]
Description=Simple log writing service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/log-systemd.sh
